# hestia-sat/systemd/uts-web-monitor.service
[Unit]
Description=UTS web monitor

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart uts-web

[Install]
WantedBy=multi-user.target
