# hestia-sat/systemd/uts-web-nginx.path
[Path]
PathChanged=/home/debian/uts/nginx/uts-web-nginx.conf

[Install]
WantedBy=multi-user.target
