# hestia-sat/systemd/uts-log-monitor.service
[Unit]
Description=UTS log monitor

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart uts-log

[Install]
WantedBy=multi-user.target
