# hestia-sat/systemd/uts-log-monitor.path
[Path]
PathChanged=/home/debian/uts/bin/uts-log

[Install]
WantedBy=multi-user.target
