# hestia-sat/systemd/uts-web-monitor.path
[Path]
PathChanged=/home/debian/uts/bin/uts-web

[Install]
WantedBy=multi-user.target
